# configs/supervisord.d/run.conf
[program:run]
command=/opt/app-root/bin/run.sh
stdout_logfile=/proc/1/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
stopasgroup=true
